# supervisord.conf
[supervisord]
nodaemon=true

[program:plone]
command=/opt/iuem.sequestre/bin/instance console
autostart=true
autrestart=true
user=plone
